# broker/etc/init/ave-broker.conf
description "ave-broker"

start on (filesystem and net-device-up IFACE!=lo)

console output
expect daemon

exec /usr/bin/ave-broker --start --force
