# handset/etc/init/ave-adb-server.conf
description "ave-adb-server"

start on (filesystem and net-device-up IFACE!=lo)

exec /usr/bin/ave-adb-server --start --force
